# bambuddy-clearplate.service
[Unit]
Description=BamBuddy Clear Plate Controller
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Bambuddy_Touch
ExecStart=/usr/bin/python3 /home/pi/Bambuddy_Touch/backend.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
